# inventory.service
[Unit]
Description=Auto Start inventory.py

[Service]
Type=simple
ExecStart=/inventory/inventory.py

[Install]
WantedBy=multi-user.target
